# linux/make-package/template-daemon/deb-scripts/postinst
#!/bin/sh
DAEMON_NAME=mydaemon
INITD=/etc/init.d/mydaemon

if [ -x /sbin/insserv ]; then
  insserv ${DAEMON_NAME}
else
  update-rc.d ${DAEMON_NAME} defaults 90
fi

$INITD restart

exit 0
