# linux/debian/php7.4-compat/package-files/install.sh
#!/bin/sh
ldconfig
systemctl enable php7.4-fpm.service
systemctl start  php7.4-fpm.service