# linux/make-package/template-daemon/data/event.d.conf.in
start on runlevel [23]
stop on runlevel [!23]

respawn
exec /usr/local/bin/@PKGNAME@