# linux/make-package/template-daemon/data/init.d.conf.in
start on (local-filesystems
          and started dbus)
stop on runlevel [!23]

respawn
exec /usr/local/bin/@PKGNAME@
